# hamonize-agent/debian/prerm
#!/bin/bash
systemctl stop hamonize-agent
systemctl disable hamonize-agent
#systemctl daemon-reload
sudo rm -r /lib/systemd/system/hamonize-agent.service

systemctl daemon-reload

# remove device contol
# cd /usr/share/hamonize-agent/usb-lockdown
# make clean

# delete hamonize-agent folder
rm -fr /usr/share/hamonize-agent
